# bin/gpucontrol.service
[Unit]
Description=PowerControl
After=multi-user.target
[Service]
Type=simple
ExecStart=/usr/local/bin/gpucontrol __gpucontrol_monitor__
StandardOutput=append:/var/log/gpucontrol.log
StandardError=append:/var/log/gpucontrol.log
[Install]
WantedBy=multi-user.target
