# bin/no_turbo.conf
# /etc/init/no_turbo.conf
description "Disable Intel Turbo Boost at startup"
start on startup
task
script
    echo 1 > /sys/devices/system/cpu/intel_pstate/no_turbo
end script
